# homeserver.jcgs.crontab
0 3 * * * /mnt/usbmem/JCGS-scripts/pullrepos
